# peterad.socket.in
[Unit]
Description=peterad socket

[Socket]
ListenStream=@PETERA_SOCKET@

[Install]
WantedBy=multi-user.target
