# peterad.socket
[Unit]
Description=peterad socket

[Socket]
ListenStream=5700

[Install]
WantedBy=multi-user.target
