# petera-decryptd.socket.in
[Unit]
Description=Petera decryptd socket

[Socket]
ListenStream=@PETERA_SOCKET@

[Install]
WantedBy=multi-user.target
