# peterad.service.in
[Unit]
Description=peterad service

[Service]
ExecStart=@PETERAD_BINARY@ @PETERA_CONF@/@PACKAGE_NAME@d.pem @PETERA_CONF@/encrypt.pem @PETERA_CONF@/decrypt.d
StandardError=syslog
