# src/plugins/askpassd/deo-askpassd.service.in
[Unit]
Description=Deo askpass daemon
DefaultDependencies=false
ConditionDirectoryNotEmpty=@DEO_CONF@/disks.d

[Service]
ExecStart=@bindir@/deo askpassd
StandardError=syslog
