# src/plugins/askpassd/petera-askpassd.service.in
[Unit]
Description=Petera askpass daemon
ConditionDirectoryNotEmpty=@PETERA_CONF@/disks.d

[Service]
ExecStart=@bindir@/petera askpassd
StandardError=syslog
