# src/plugins/askpassd/petera-askpassd.path.in
[Unit]
Description=Petera askpass daemon path

[Path]
DirectoryNotEmpty=/run/systemd/ask-password

[Install]
WantedBy=basic.target
