# src/plugins/decryptd/petera-decryptd.socket.in
[Unit]
Description=Petera decryption daemon socket

[Socket]
ListenStream=@PETERA_SOCKET@

[Install]
WantedBy=multi-user.target
