# src/plugins/decryptd/petera-decryptd.socket.in
[Unit]
Description=Petera decryption daemon socket
ConditionDirectoryNotEmpty=@PETERA_CONF@/decrypt.d
ConditionFileNotEmpty=@PETERA_CONF@/decryptd.pem
ConditionPathExists=@PETERA_CONF@/encrypt.pem

[Socket]
ListenStream=@PETERA_SOCKET@

[Install]
WantedBy=multi-user.target
